# systemd/vpn-remote-bot.service
[Unit]
Description=VPN Remote Updater Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
Environment="PYTHONUNBUFFERED=1"
ExecStart=__VENV_DIR__/bin/python -u bot/bot.py
Restart=on-failure
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target